# configs/stash-cache/systemd/stash-cache-authfile.service
[Unit]
Description=StashCache Authfile generator
After=network-online.target

[Service]
User=xrootd
Group=xrootd
Type=oneshot
ExecStart=/usr/libexec/xcache/authfile-update --cache

[Install]
WantedBy=multi-user.target
