# configs/stash-cache-authfile.service
[Unit]
Description=StashCache Authfile generator

[Service]
User=xrootd
Group=xrootd
Type=oneshot
ExecStart=/usr/libexec/xcache/authfile-update

[Install]
WantedBy=multi-user.target
