# configs/stashcache-authfile.service
[Unit]
Description=StashCache Authfile generator

[Service]
User=xrootd
Group=xrootd
Type=oneshot
ExecStart=/usr/libexec/stashcache-cache-server-auth/authfile-update

[Install]
WantedBy=multi-user.target
