# configs/stashcache-authfile-public.service
[Unit]
Description=StashCache public Authfile generator

[Service]
User=xrootd
Group=xrootd
Type=simple
ExecStart=/usr/libexec/stashcache-cache-server/authfile-public-update

[Install]
WantedBy=multi-user.target
